# idioms-quial/src/main/kotlin/com/backend/quial/python/crontab
30 22 * * * root /usr/bin/python3 /app/quial.py >> /var/log/cron.log; root /usr/bin/python3 /app/topic-quial.py >> /var/log/cron2.log
