# idioms-quial/src/main/kotlin/com/backend/quial/python/crontab
30 23 * * * root /usr/bin/python3 /app/quial.py >> /var/log/cron.log
